# maintenance/config/systemd/happyserv-auto-update.service
[Unit]
Description=HappyServ Auto-Update Service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/home/happyserv/maintenance/scripts/auto-update.sh
WorkingDirectory=/home/happyserv/maintenance
StandardOutput=journal
StandardError=journal
