# maintenance/config/systemd/happyserv-health-check.service
[Unit]
Description=HappyServ Health Check Service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/home/happyserv/maintenance/scripts/health-check.sh
WorkingDirectory=/home/happyserv/maintenance
StandardOutput=journal
StandardError=journal
