# config/systemd/happyserv-redis.service
[Unit]
Description=HappyServ Redis Service
After=docker.service
Requires=docker.service
BindsTo=happyserv-docker.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f /home/happyserv/docker/compose/prod/docker-compose.yml up -d redis
ExecStop=/usr/bin/docker compose -f /home/happyserv/docker/compose/prod/docker-compose.yml stop redis
WorkingDirectory=/home/happyserv
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=happyserv-docker.target
