# config/systemd/happyserv-api.service
[Unit]
Description=HappyServ API Service
After=docker.service
Requires=docker.service
BindsTo=happyserv-docker.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -f /home/happyserv/docker/compose/prod/docker-compose.yml up -d api
ExecStop=/usr/bin/docker compose -f /home/happyserv/docker/compose/prod/docker-compose.yml stop api
ExecReload=/usr/bin/docker compose -f /home/happyserv/docker/compose/prod/docker-compose.yml restart api
WorkingDirectory=/home/happyserv
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=happyserv-docker.target
